# deploy/systemd/gsad-gpu-server-report.service
[Unit]
Description=GSAD GPU metrics reporter
Documentation=file:@REPO_ROOT@/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@REPO_ROOT@/gpu-server-report
EnvironmentFile=-@REPO_ROOT@/deploy/env/common.env
EnvironmentFile=-@REPO_ROOT@/deploy/env/reporter.env
ExecStart=@UV_BIN@ run python reporter.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
